# supervisor.ini
[program:wechat]
command=python3 main.py              ; the program (relative uses PATH, can take args)
directory=/root/code/zju_wechat               ; directory to cwd to before exec (def no cwd)
priority=999                  ; the relative start priority (default 999)
autostart=true                ; start at supervisord start (default: true)
autorestart=true              ; retstart at unexpected quit (default: true)
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/root/code/zju_wechat/supervisor.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/root/code/zju_wechat/supervisor.log        ; stderr log path, NONE for none; default AUTO

[program:tem]
command=python3 main.py              ; the program (relative uses PATH, can take args)
directory=/root/code/zju_tem               ; directory to cwd to before exec (def no cwd)
priority=999                  ; the relative start priority (default 999)
autostart=true                ; start at supervisord start (default: true)
autorestart=true              ; retstart at unexpected quit (default: true)
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/root/code/zju_tem/supervisor.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/root/code/zju_tem/supervisor.log        ; stderr log path, NONE for none; default AUTO

